# etc/supervisord.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-start
